# packages/core/host/networking/mac-addresses/overlays/usr/lib/systemd/system/report-mac-addresses.timer
[Unit]
Description=Run the report-mac-addresses service periodically in case the connected network interfaces change

[Timer]
OnUnitActiveSec=2min
OnUnitInactiveSec=1min

[Install]
WantedBy=timers.target
