# packages/core/host/boot-init-root/overlays/usr/lib/systemd/system/boot-init-root.service
[Unit]
Description=Move files from `/boot/init-root` into the root filesystem
ConditionPathExistsGlob=/boot/init-root/*

[Service]
Type=oneshot
ExecStart=sh -c 'mv /boot/init-root/* /'
Restart=no

[Install]
WantedBy=multi-user.target
