# packages/core/host/boot-init-root/overlays/usr/lib/systemd/system/boot-init-root.service
[Unit]
Description=Move files from /boot/firmware/init-root into the root filesystem
ConditionPathExistsGlob=/boot/firmware/init-root/*

[Service]
Type=oneshot
ExecStart=bash -c '/usr/libexec/boot-init-root'
Restart=no

[Install]
WantedBy=multi-user.target
