# deployments/provisioning/boot-init-root.pkg/overlays/usr/lib/systemd/system/boot-init-root-extract-archives.service
[Unit]
Description=Extract init-root archives in /boot/firmware to /boot/firmware/init-root
ConditionPathExistsGlob=/boot/firmware/init-root-*.tar.gz

[Service]
Type=oneshot
ExecStart=bash -c '/usr/libexec/boot-init-root-extract-archives'
Restart=no

[Install]
WantedBy=multi-user.target
